# etc/profile.d/aliases.sh
################################################################################
##
## Alces Metalware - Shell configuration
## Copyright (c) 2008-2015 Alces Software Ltd
##
################################################################################
alias power="metal power"
alias console="metal console"
